# kb-layout-daemon.service
[Unit]
Description=Keyboard Layout Daemon - Auto-switch layout per keyboard
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/kb-layout-daemon
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info

[Install]
WantedBy=default.target
